# backend/utils/supervisord.conf
[program:TartuNightLife]
command=src/server.py
